# Plugins/EASy-Producer/ScenariosTest/testdata/real/IIP-Ecosphere/experiment/expected/SerializerConfig1/iip-ecs.service
[Unit]
Description=IIP-Ecosphere ECSRuntime
After=iip-platform.service
Requires=iip-platform.service

[Service]
WorkingDirectory=/opt/iip/gen
Type=simple
Restart=on-failure
KillMode=control-group
PIDFile=/run/iip-ecsRuntime.pid
ExecStart=/opt/iip/java -cp @ecsJars/classpath.ux -Dio.netty.tryReflectionSetAccessible=true -Dlog4j2.formatMsgNoLookups=true -Xmx512M --add-opens java.base/jdk.internal.misc=ALL-UNNAMED --add-opens java.base/java.lang=ALL-UNNAMED -Diip.pid.dir=/run "de.iip_ecosphere.platform.support.LifecycleHandler$WaitingStarter"

[Install]
WantedBy=multi-user.target
